# config/macos/com.investigator.readyshare.mount.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.investigator.readyshare.mount</string>

    <key>ProgramArguments</key>
    <array>
      <string>/usr/bin/env</string>
      <string>bash</string>
      <string>/ABSOLUTE/PATH/TO/InvestiGator/scripts/mount_readyshare.sh</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
      <key>READYSHARE_CREDENTIALS</key>
      <string>/etc/readyshare-credentials</string>
      <key>READYSHARE_HOST</key>
      <string>192.168.1.1</string>
      <key>READYSHARE_SHARE</key>
      <string>USB_Storage</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>/ABSOLUTE/PATH/TO/InvestiGator</string>

    <key>RunAtLoad</key>
    <true/>
    <key>StartInterval</key>
    <integer>300</integer>

    <key>KeepAlive</key>
    <dict>
      <key>SuccessfulExit</key>
      <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/readyshare-mount.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/readyshare-mount.log</string>
  </dict>
</plist>
